# docker/playwright-remote/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1440x900x24 -ac +extension RANDR
autorestart=true
priority=10

[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":99"
autorestart=true
priority=20

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -shared -nopw -xkb
autorestart=true
priority=30

[program:novnc]
command=/usr/share/novnc/utils/novnc_proxy --vnc localhost:5900 --listen 6080
autorestart=true
priority=40

[program:chromium]
command=/usr/local/bin/start-chromium.sh
environment=DISPLAY=":99"
autorestart=true
priority=50
